# crontab
*/10 * * * * root hdparm -Y /dev/sdb
*/10 * * * * root hdparm -Y /dev/sdc
